# tools/linux/d06-evdev.service
[Unit]
Description=D06 Pro evdev decoder
Documentation=https://github.com/SupperTomato/d06-pro-ring-mouse-toolkit

[Service]
Type=simple
WorkingDirectory=%h/d06-pro-ring-mouse-toolkit
ExecStart=/usr/bin/python3 %h/d06-pro-ring-mouse-toolkit/tools/linux/d06_evdev.py --profile %h/.config/d06/profile.json
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
